# crontab
# crontab
# นาที ชั่วโมง วัน เดือน วันในสัปดาห์ คำสั่ง

# ส่งแจ้งเตือนออกงานเวลา 17:30 (5:30 PM) Asia/Bangkok วันจันทร์ถึงศุกร์
# เวลานี้เหมาะสมสำหรับการแจ้งเตือนก่อนเลิกงาน
30 17 * * 1-5 curl -H "Authorization: Bearer $CRON_SECRET" http://app:12914/api/cron/checkout-reminder

# ลงชื่อออกงานอัตโนมัติตอนเที่ยงคืน (00:00) สำหรับคนที่ลืมลงชื่อออก
# ทำงานทุกวัน เพื่อให้ครอบคลุมทุกวันทำงาน
0 0 * * * curl -H "Authorization: Bearer $CRON_SECRET" http://app:12914/api/cron/auto-checkout

# ตรวจสอบให้แน่ใจว่า container และ host มี timezone ที่ถูกต้อง
# ถ้า Dockerfile.cron ตั้ง TZ=Asia/Bangkok แล้ว schedule นี้จะทำงานตามเวลา Asia/Bangkok

# บรรทัดสุดท้ายต้องเป็น newline เสมอสำหรับ crontab
